# deploy/gcp-l4/nauticai-api.service
[Unit]
Description=NautiCAI FastAPI (GPU inference)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/nauticai/app
Environment=PYTHONPATH=/opt/nauticai/app
EnvironmentFile=/opt/nauticai/app/backend/.env
ExecStart=/opt/nauticai/app/.venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=on-failure
RestartSec=5
# Model warmup can take 1–2 minutes on cold start
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
